# templates/rabbitmq-server.service.d/limits.conf
[Service]
LimitNOFILE=<%= @file_limit %>
Restart=<%= @restart_param %>
StartLimitInterval=<%= @start_limit_interval %>
StartLimitBurst=<%= @start_limit_burst %>
